# deploy/amyzhou-leads.service
[Unit]
Description=Amy Zhou Homes private lead intake API
After=network.target

[Service]
Type=simple
DynamicUser=yes
StateDirectory=amyzhou-leads
StateDirectoryMode=0700
EnvironmentFile=-/etc/amyzhou-leads.env
ExecStart=/usr/bin/python3 /opt/amyzhou-lead-api/lead_api.py
Restart=on-failure
RestartSec=3
UMask=0077
NoNewPrivileges=yes
ProtectProc=invisible
ProcSubset=pid
PrivateDevices=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
CapabilityBoundingSet=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
